# files/etc/rc.local
# Put your custom commands here that should be executed once
#sh /usr/bin/jam_oc.sh vh2.onelink.me
#sh /usr/bin/jam_nikki.sh vh2.onelink.me
exit 0
